# root/etc/init.d/lan-proxy-route
#!/bin/sh /etc/rc.common

START=95
STOP=05
USE_PROCD=1

SERVICE=/usr/share/lan-proxy-route/lan-proxy-route.sh

start_service() {
	"$SERVICE" diagnose >/dev/null
	"$SERVICE" apply
}

reload_service() {
	"$SERVICE" cleanup
	"$SERVICE" diagnose >/dev/null
	"$SERVICE" apply
	/etc/init.d/dnsmasq reload >/dev/null 2>&1 || true
	/etc/init.d/firewall reload >/dev/null 2>&1 || true
}

stop_service() {
	"$SERVICE" cleanup
}
